# systemd/adblocker.service
[Unit]
Description=eBPF Adblocker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/adblockerctl daemon --config /etc/adblocker/adblocker.yaml
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_BPF CAP_NET_ADMIN CAP_SYS_RESOURCE
CapabilityBoundingSet=CAP_BPF CAP_NET_ADMIN CAP_SYS_RESOURCE
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/sys/fs/bpf
ProtectHome=true

[Install]
WantedBy=multi-user.target
